# baseboard-rpi3/chromeos-base/snd_bcm2835-spec/files/force_audio_output_to_headphones.conf
# Copyright (c) 2019 The FydeOS Authors. All rights reserved.
# Use of this source code is governed by a BSD-style license that can be
# found in the LICENSE file.

description     "Force setting audio output to headphone for Raspberry Pi 3b/3b+/4b"
author          "yang@fydeos.io"

# The cras (chrome os audio service) is uncompatiable with Pi's audio driver, the 
# default output is "auto", that will make cras confused.
# According to https://elinux.org/R-Pi_Troubleshooting#Sound
# If you want to enable hdmi audio, modify the last param to 
#         0    # unplug headphone, use hdmi instead (if hdmi audio is exist)
#         1    # plug headphone

start on starting login

env CRAS_CLIENT=/usr/bin/cras_test_client

script
  dev_num=$($CRAS_CLIENT --dump_s | grep HEADPHONE | awk '{print $2}')
  $CRAS_CLIENT --plug ${dev_num}:1
end script
