# project-cros-pi/chromeos-base/vpd/files/check_serial_number.conf
description   "Check serial number for compability"
author        "hi@fydeos.io"

start on starting system-services
exec /bin/sh /usr/share/cros/init/check_serial_number.sh
