# baseboard-rpi3/media-sound/adhd/files/cras_monitor.conf
# Copyright 2019 The FydeOS Authors. All rights reserved.
# Use of this source code is governed by a BSD-style license that can be
# found in the LICENSE file.

description   "Monitor cras status, fix sound stream stucked issue"
author        "yang@fydeos.io"

start on login-prompt-visible
stop on stopping boot-services

respawn
respawn limit 3 10

expect fork

exec /usr/bin/cras_monitor
