# project-cros-pi/chromeos-base/workaround-start-ui-after-vc4-inited/files/rules/99-vc4-monitor.rules
ACTION=="add", KERNEL=="card[0-9]", DRIVERS=="vc4-drm", RUN+="/sbin/start splash"
