# overlay-rpi4-fydeos/chromeos-base/chromeos-bsp-rpi4-fydeos/files/audio/rpi4-hdmi.conf
description   "Plug HDMI manually"
author        "yang@fydeos.io"

start on start-user-session
task
oom score -100

exec /bin/bash /usr/share/cros/init/set-hdmi.sh
