# overlay-rpi4/chromeos-base/chromeos-bsp-rpi4/files/audio/rpi4-hdmi.conf
description   "Plug HDMI manually"
author        "yang@fydeos.io"

start on start-user-session
task
oom score -100

exec /usr/share/cros/init/set-hdmi.sh
