# overlay-rpi3/chromeos-base/chromeos-bsp-rpi3/files/audio/force_audio_output_to_headphones.conf
# Copyright (c) 2019 The FydeOS Authors. All rights reserved.
# Use of this source code is governed by a BSD-style license that can be
# found in the LICENSE file.

description     "Force setting audio output to headphone for Raspberry Pi 3b/+"
author          "yang@fydeos.io"

# The cras (chrome os audio service) is uncompatiable with Pi's audio driver, the
# default output is "auto", that will make cras confused.
# According to https://elinux.org/R-Pi_Troubleshooting#Sound
# If you want to enable hdmi audio, modify the last number to 2

start on start-user-session

script
  /usr/bin/amixer -c 1 cset numid=3 1
  /usr/bin/cras_test_client --select_output 6:0
end script
