# overlay-rpi4/sys-apps/haveged/files/haveged-fydeos.conf
# Copyright (c) 2017 The Flint OS Authors. All rights reserved.
# Use of this source code is governed by a BSD-style license that can be
# found in the LICENSE file.

description     "Software RNG Entropy Gatherer Daemon"
author          "yang@flintos.io"

#start on started boot-services
start on stopped udev-trigger-early
stop on stopping pre-shutdown


script
  WATERMARK=1024

# -r0 is added always
  HAVEGED_OPTS="-w ${WATERMARK} -v 1"

  SVCNAME=haveged
  pidfile=/run/${SVCNAME}.pid

  command="/usr/sbin/${SVCNAME}"
  command_args="-r 0 ${HAVEGED_OPTS} -p ${pidfile}"
  exec $command $command_args
end script
