# project-fyde-for-rpi/chromeos-base/vpd/files/check_serial_number.conf
description   "Check serial number for compability"
author        "yang@fydeos.io"

start on starting system-services
exec /bin/sh /usr/share/cros/init/check_serial_number.sh
