# baseboard-rpi3/net-wireless/bluez/files/bluez-upstart.conf
# Copyright (c) 2012 The Chromium OS Authors. All rights reserved.
# Use of this source code is governed by a BSD-style license that can be
# found in the LICENSE file.

description     "Start the bluetooth daemon"
author          "chromium-os-dev@chromium.org"

start on started system-services
stop on stopping system-services

limit memlock unlimited unlimited

env BLUETOOTH_LIBDIR=/var/lib/bluetooth
env BLUETOOTH_DAEMON_OPTION=""

pre-start script
  mkdir -p -m 0750 ${BLUETOOTH_LIBDIR}
  chown -R bluetooth:bluetooth ${BLUETOOTH_LIBDIR}
end script

respawn

exec /usr/bin/start_bluetoothd.sh
