# overlay-rpi3/chromeos-base/chromeos-bsp-rpi3/files/udev/99-rpi3-restart-ui-after-hdmi-connetced.rules
# Copyright (c) 2019 The FydeOS Authors. All rights reserved.
# Use of this source code is governed by a BSD-style license that can be
# found in the LICENSE file.

# restart ui after hdmi connected.

ACTION=="add", SUBSYSTEM=="drm", KERNEL=="renderD128", RUN+="/lib/udev/restart_ui.sh %k"
